# kit/etc/udev/rules.d/10-vchiq-permissions.rules
SUBSYSTEM=="vchiq",GROUP="video",MODE="0660",RUN+="/opt/vc/bin/add_shm_tree.sh"
